# etc/cron.daily/update_stats
#!/bin/sh
su - pi <<EOF
python /opt/smac/bin/update_stats.py
EOF
